# java/bin/canerozdemircgi.service
[Unit]
Description=canerozdemircgi service

[Service]
Type=simple
WorkingDirectory=/var/www/html/java/bin
ExecStart=java -jar canerozdemircgi.jar

[Install]
WantedBy=multi-user.target